# deployment/Kafka/kafka-connect.conf
[program:kafka-connect]
command=/usr/bin/start-connect.sh
autostart=true
autorestart=true